# humhub-docker/base/etc/supervisord.conf.d/cron.conf
[program:cron]
command=crond -f -L /proc/self/fd/2
user=root
autostart=true
autorestart=true
