# humhub-docker/nginx/etc/supervisord.conf.d/nginx.conf
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
user=root
autostart = true
stdout_logfile=/proc/self/fd/2
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0
